# common/etc/udev/rules.d/90-backlight-perms.rules
ACTION=="add", SUBSYSTEM=="backlight", RUN+="/bin/chmod 0777 /sys/class/backlight/%k/brightness"
