# common/etc/udev/rules.d/91-cpu-perms.rules
ACTION=="add", SUBSYSTEM=="cpu", KERNEL=="cpu[0-9]*", RUN+="/bin/chmod 0777 /sys/devices/system/cpu/%k/cpufreq/scaling_governor"
