# common/etc/udev/rules.d/92-rockchip-ebc.rules
ACTION=="add", KERNEL=="rockchip_ebc_blit_neon", RUN+="/bin/chmod 777 /sys/module/rockchip_ebc_blit_neon/parameters/y4_threshold_y1"
